# workspace/dashboard/com.helm.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.helm.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/owner/.local/bin/node</string>
    <string>/Users/owner/secondme/workspace/dashboard/server.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/owner/secondme</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/owner/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/owner/secondme/workspace/dashboard/dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/owner/secondme/workspace/dashboard/dashboard.log</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
